# dossier-industrial-alimentario/scripts/systemd/surus-agente-financieros.service
[Unit]
Description=HERMES Agente Financieros (C.2) — Scrape Wikipedia para KPIs A&B
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/hermes-dossier/apps/dossier-industrial
EnvironmentFile=/opt/hermes-dossier/.env
ExecStart=/usr/bin/pnpm scan:financials
StandardOutput=append:/opt/hermes-dossier/logs/financials.log
StandardError=append:/opt/hermes-dossier/logs/financials.err
