# dossier-industrial-alimentario/scripts/systemd/surus-agente-financieros.timer
[Unit]
Description=HERMES Agente Financieros (C.2) — Timer semanal

[Timer]
OnCalendar=weekly
Persistent=true
Unit=surus-agente-financieros.service

[Install]
WantedBy=timers.target
