# dossier-industrial-alimentario/deploy/hermes-scan-linkedin.service
[Unit]
Description=HERMES Dossier — LinkedIn OSINT agent
After=network.target postgresql.service hermes-dossier.service
Wants=postgresql.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/hermes-dossier/apps/dossier-industrial
EnvironmentFile=/opt/hermes-dossier/.env
Environment=NODE_ENV=production
Environment=MAX_QUERIES=15
ExecStart=/usr/bin/node node_modules/tsx/dist/cli.mjs scripts/scan-linkedin.ts
StandardOutput=append:/opt/hermes-dossier/logs/agent-linkedin.log
StandardError=append:/opt/hermes-dossier/logs/agent-linkedin.err
TimeoutStartSec=1800s
Nice=10

[Install]
WantedBy=multi-user.target
